# systemd/containers/acme-register.container
[Unit]
Description=The acme-register container
Requires=acme-certs-volume.service openchami-cert-trust.service acme-certs-volume.service
After=acme-certs-volume.service openchami-cert-trust.service acme-certs-volume.service
PartOf=openchami.target

[Container]
ContainerName=demo.openchami.cluster
HostName=demo.openchami.cluster
Image=docker.io/neilpang/acme.sh:3.1.1

# Volumes
Volume=/etc/pki/ca-trust/extracted/pem/tls-ca-bundle.pem:/root_ca/root_ca.crt:ro,Z
Volume=acme-certs:/acme.sh:z

# Networks for the Container to use
Network=openchami-cert-internal.network

Exec=--issue \
  --ca-bundle /root_ca/root_ca.crt \
  --server https://step-ca:9000/acme/acme/directory \
  --home /acme.sh \
  -d demo.openchami.cluster \
  --standalone \
  --force

[Service]
Restart=on-failure
Type=oneshot
RemainAfterExit=yes
