# systemd/containers/opaal-idp.container
[Unit]
Description=The opaal-idp container
Requires=hydra-gen-jwks.service
After=hydra-gen-jwks.service

[Container]
ContainerName=opaal-idp
HostName=opaal-idp
Image=ghcr.io/openchami/opaal:v0.3.10

# Volumes
Volume=/etc/openchami/configs:/opaal/config/

# Command to run in container
Exec=/opaal/opaal serve --config /opaal/config/opaal.yaml

# Networks for the Container to use
Network=openchami-jwt-internal

# Proxy settings
PodmanArgs=--http-proxy=false

[Service]
Restart=always
